# dns/supervisord.conf
[supervisord]
nodaemon=true
pidfile=/var/run/supervisord.pid
logfile=/var/log/supervisord.log

[unix_http_server]
file=/var/run/supervisor.sock
chmod=0700

[program:bind9]
command=/usr/sbin/named -f -c /etc/bind/named.conf -u root
autorestart=true
user=root

[program:bind9-api]
command=/usr/local/bin/bind9-api --config /etc/bind9-api/config.yaml
directory=/etc/bind9-api
autorestart=true
user=root